# debian/postinst
#!/bin/bash

# Install daemon
cp /usr/share/gooroom/security/os-protector/gop-daemon /etc/init.d/
update-rc.d gop-daemon defaults

# Disable interrupt remapping to reserve IOMMU if Gooroom Boot Protector is not installed
service gbp-daemon status > /dev/null
ret=$?

# Run update-grub if gooroom-boot-protector is not installed
dpkg -l | grep -q gooroom-boot-protector
if [ $? -ne 0 ]; then
    if [ \( $ret -ne 0 \) -a \( $ret -ne 3 \) ]; then
    	sed -i 's/GRUB_CMDLINE_LINUX_DEFAULT=\"/GRUB_CMDLINE_LINUX_DEFAULT=\"reboot=w,a,f iommu=force intel_iommu=on /' /etc/default/grub
    	update-grub
    fi
fi

## Delete previous plka file
#rm -f /etc/polkit-1/localauthority/90-mandatory.d/disable-suspend.pkla
#rm -f /etc/polkit-1/localauthority/90-mandatory.d/disable-hibernate.pkla

# Ignore power keys for avoiding Samsung Notebook Issue
sed -i 's/#HandleSuspendKey=suspend/HandleSuspendKey=ignore/; s/#HandleHibernateKey=hibernate/HandleHibernateKey=ignore/; s/#HandleLidSwitch=suspend/HandleLidSwitch=ignore/; s/#HandlePowerKey=poweroff/HandlePowerKey=ignore/' /etc/systemd/logind.conf
sed -i 's/#xserver-command=X/xserver-command=X -s 0 dpms/' /etc/lightdm/lightdm.conf
