# debian/prerm
#!/bin/bash

# Uninstall daemon
update-rc.d os_protection remove
rm /etc/init.d/os_protection
