# debian/prerm
#!/bin/bash

# Uninstall daemon
update-rc.d gop-daemon remove
rm /etc/init.d/gop-daemon
